# py-human-tracking.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/my-python-script.py

[Install]
WantedBy=default.target